# roles/nfs-client/templates/share.automount
[Unit]
Description=Automount share

[Automount]
Where={{ nfs_share_mountpoint }}

[Install]
WantedBy=multi-user.target
